# sr_auth.service.template
[Unit]
Description=SR auth server service
[Service]
Type=simple
Restart=always
ExecStart=SERVER_CMD
User=SERVER_USER
[Install]
WantedBy=multi-user.target
